# p-cast.service
[Unit]
Description=p-cast: stream PipeWire audio to Chromecast
After=network.target pipewire.service pipewire-pulse.service
ConditionUser=!root

# Start with `systemctl --user start p-cast`
#
# Type=notify — p-cast signals to systemd once startup is complete; systemd tracks the PID automatically.
# StateDirectory=p-cast — creates ~/.local/state/p-cast/ automatically (for the log file); persistent across reboots.
# Restart=on-failure / RestartSec=10 — restarts if p-cast crashes or if Chromecast discovery fails at boot before the devices are online.
# WantedBy=default.target — correct target for user services (vs multi-user.target for system services).
# TimeoutStopSec=15 — gives p-cast 15 seconds to stop cleanly before systemd force-kills it

[Service]
Type=notify
NotifyAccess=all
StateDirectory=p-cast
ExecStart=%h/.local/bin/p-cast \
    --log-file %h/.local/state/p-cast/p-cast.log

Restart=on-failure
RestartSec=10
TimeoutStopSec=15

[Install]
WantedBy=default.target
